# features/gardener/file.include/etc/systemd/system/overlay.mount
[Mount]
What=none
Where=/overlay
Type=tmpfs
Options=mode=0755

[Install]
WantedBy=local-fs.target
