# features/_pxe/file.include/usr/lib/dracut/modules.d/98gardenlinux-live/live-get-squashfs.service
[Unit]
Description=Download squashfs image

After=network-online.target 
Wants=network-online.target

OnFailure=emergency.target
OnFailureJobMode=isolate
DefaultDependencies=no

[Service]
Type=oneshot
TimeoutStartSec=600
RemainAfterExit=yes
ExecStart=/sbin/live-get-squashfs
