# features/azure/file.include/etc/systemd/system/chronyd-restricted.service.d/05-private-devices-access.conf
[Service]
DeviceAllow=/dev/ptp_hyperv
DevicePolicy=closed
PrivateDevices=no
