# features/azure/file.include/etc/systemd/system/walinuxagent.service.d/10-startpost.conf
[Service]
ExecStartPost=/bin/systemctl try-restart systemd-networkd.service
