# features/_usi/file.include/etc/systemd/system/reload.target
[Unit]
After=local-fs.target
Wants=local-fs.target
Wants=systemd-repart.service
Wants=systemd-udevd.service systemd-udev-trigger.service
Wants=systemd-journald.service
Requires=reload.service
