# features/_usi/file.include/etc/systemd/system/root.mount
[Mount]
What=none
Where=/root
Type=tmpfs
Options=size=128M,mode=0700

[Install]
WantedBy=reload.target
