# features/_nocrypt/initrd.include/etc/systemd/system/sysroot-var.mount
[Unit]
Requires=systemd-repart.service
After=systemd-repart.service

[Mount]
What=/dev/disk/by-partlabel/VAR
Where=/sysroot/var
Type=ext4
