# features/server/file.include/etc/update-motd.d/10-hostname
#!/bin/sh
echo "Welcome to \e[0;32m$(hostname -f)\e[0m  $(date)"
