# features/server/file.include/etc/update-motd.d/20-uname
#!/bin/sh
echo "$(uname -svm) $(uptime -p)"
