# features/aws/file.include/etc/systemd/system/cloud-init-local.service.d/override.conf
# temporary fix so that the name resolution doesn't block
# TODO: fix this in a nicer way
[Unit]
Before=systemd-resolved-varlink.socket
