# features/server/delete.list
/etc/machine-id
/var/log/fail2ban.log
/var/log/journal
/etc/init.d/sysstat
/etc/ufw
/etc/init.d/auditd
/etc/init.d/dbus
/etc/init.d/fail2ban
/etc/monit
/etc/init.d/ferm
/etc/init.d/ssh
/etc/init.d/sudo
/etc/init.d/sysstat
/etc/init.d/udev
