# features/server/file.include/etc/systemd/system/systemd-timesyncd.service.d/override.conf
[Unit]
After=systemd-resolved.service
Wants=systemd-resolved.service
