# features/_dev/file.include/etc/systemd/system/emergency.service.d/sulogin.conf
[Service]
Environment=SYSTEMD_SULOGIN_FORCE=1
