# features/server/file.include/etc/systemd/system/systemd-coredump@.service.d/override.conf
[Service]
ExecStopPost=/bin/sh -c "rmdir /run/systemd/propagate/systemd-coredump@%i.service"
