# features/_tpm2/initrd.include/etc/systemd/system/check-tpm.service
[Unit]
DefaultDependencies=no
After=tpm2.target
Wants=tpm2.target
Before=systemd-repart.service

[Service]
Type=oneshot
ExecStart=/usr/bin/check-tpm
RemainAfterExit=yes
