# features/_tpm2/initrd.include/etc/systemd/system/sysroot-var.mount
[Unit]
Wants=systemd-cryptsetup-var.service
After=systemd-cryptsetup-var.service

[Mount]
What=/dev/mapper/var
Where=/sysroot/var
Type=ext4
