# features/_kdump/file.include/etc/systemd/system/kdump-tools.service.d/override.conf
[Service]
ExecStartPre=/usr/local/sbin/prepare-initrd-kdump
