# features/_debug/file.include/etc/systemd/system/setup-debugbox-ssh.service
[Unit]
Description=SSH Container Setup
PropagatesReloadTo=ssh.service

[Service]
#RemainAfterExit=true
Type=oneshot
ExecStart=/usr/sbin/setup-debugbox-ssh

[Install]
WantedBy=ssh.service