# features/_trustedboot/file.include/etc/systemd/system/tmpfsify-etc-systemd.service
[Unit]
Before=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tmpfsify /etc/systemd/system

[Install]
WantedBy=local-fs.target
