# features/_trustedboot/file.include/etc/systemd/system/var.mount
[Mount]
What=none
Where=/var
Type=tmpfs
Options=mode=0755

[Install]
WantedBy=local-fs.target
