# features/_ephemeral/initrd.include/etc/systemd/system/sysroot-var.mount
[Unit]
Wants=ephemeral-cryptsetup.service
After=ephemeral-cryptsetup.service

[Mount]
What=/dev/mapper/ephemeral
Where=/sysroot/var
Type=ext4
