# features/_usi/initrd.include/etc/systemd/system/sysroot-home.mount
[Mount]
What=/sysroot/var/home
Where=/sysroot/home
Type=none
Options=bind
