# features/_usi/initrd.include/etc/systemd/system/sysroot-root.mount
[Mount]
What=tmpfs
Where=/sysroot/root
Type=tmpfs
Options=size=128M,mode=0700
