# features/_usi/initrd.include/etc/systemd/system/sysroot-opt.mount
[Mount]
What=/sysroot/var/opt
Where=/sysroot/opt
Type=none
Options=bind
