# features/_usi/initrd.include/etc/systemd/system/sysroot.mount
[Mount]
What=/root.img
Where=/sysroot
Type=erofs
