# features/_tpm2/file.include/etc/systemd/system/check-tpm.service
[Unit]
DefaultDependencies=no
After=tpm2.target
Before=local-fs-pre.target
Before=systemd-repart.service
OnFailure=emergency.target

[Service]
Type=oneshot
ExecStart=/usr/bin/check-tpm
RemainAfterExit=yes

[Install]
RequiredBy=local-fs.target
