# features/azure/file.include/etc/systemd/system/chronyd.service.d/10-after_dev-ptp_hyperv.device.conf
[Unit]
Requires=dev-ptp_hyperv.device
After=dev-ptp_hyperv.device
