# features/openstackMetal/file.include/etc/profile.d/50-autologout.sh
# shellcheck shell=sh disable=SC2148
TMOUT=900
readonly TMOUT
export TMOUT
