# features/server/copy/etc/systemd/system/systemd-resolved.service.d/resolv.conf
[Service]
ExecStartPost=+/bin/rm /etc/resolv.conf
ExecStartPost=+/bin/ln -sf /run/systemd/resolve/resolv.conf /etc/resolv.conf
