# features/disaSTIGmedium/file.include/etc/systemd/system/aide-check.timer
[Unit]
Description=Run AIDE integrity check daily at 10:00

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target

