# features/metal/file.include/etc/udev/rules.d/69-nostbyrot.rules
ACTION=="add|change", SUBSYSTEM=="block", ATTRS{queue/rotational}=="1", RUN+="/usr/bin/hdparm -B 254 -S 0 /dev/%k"
