# features/metal/delete.list
/etc/network
/etc/init.d/haveged
/etc/init.d/ipmievd
/etc/runit
/etc/sv
/etc/init.d/irqbalance
/usr/share/runit
