# features/fedramp/file.include/var/spool/cron/crontabs/root
0 30 * * * root /usr/bin/clamscan -ri / >> /var/log/clamav/clamd.log
