# features/server/systemd-hostnamed.service.d/persist.conf
[Service]
ExecStartPost=bash -c "hostnamectl set-hostname $(hostnamectl --transient)"
