# features/server/test/debsums_exclude.list
/lib/systemd/system/dbus.socket
/usr/share/pam-configs/cracklib
