# features/_immutable/file.include/usr/lib/systemd/system/systemd-tmpfiles-setup.service.d/sysroot.conf
[Unit]
After=-.mount sysroot.mount systemd-sysusers.service systemd-journald.service