# features/aws/file.include/etc/systemd/system/aws-clocksource.service
[Unit]
Description=AWS Clocksource Setup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/clocksource-setup.sh

[Install]
WantedBy=basic.target
