# features/cloud/copy/etc/udev/rules.d/75-cloud-ifupdown.rules
# Handle allow-hotplug interfaces
SUBSYSTEM=="net", ACTION=="add", ENV{INTERFACE}=="eth*|en*", RUN+="/etc/network/cloud-ifupdown-helper"
