# features/github_action_runner/file.include/etc/systemd/system/configure_from_aws_user_data.service
[Unit]
After=network-online.target

[Service]
ExecStart=/usr/bin/configure_from_aws_user_data

[Install]
WantedBy=multi-user.target
