# features/ghost/file.include/etc/systemd/system/docker.service.d/10-makeitstarted.conf
[Service]
ExecStartPost=/usr/bin/docker ps

