# features/_usi/file.include/etc/systemd/system/auditd.service.d/mkdir.conf
[Service]
ExecStartPre=/usr/bin/mkdir -p /var/log/audit
