# features/server/file.include/etc/profile.d/50-nohistory.sh
# shellcheck shell=sh disable=SC2148
HISTFILE=/dev/null
readonly HISTFILE
export HISTFILE
