# features/server/file.include/etc/profile.d/50-nohistory.sh
HISTFILE=/dev/null
readonly HISTFILE
export HISTFILE
