# features/gardener/file.include/etc/systemd/system/containerd.service.d/override.conf
[Service]
LimitMEMLOCK=67108864
LimitNOFILE=1048576
