# features/aide/file.include/var/spool/cron/crontabs/root
0 10 * * * root /usr/bin/aide --check --config /etc/aide/aide.conf > /dev/shm/aide_check 2>&1
