# features/glvd/file.include/etc/update-motd.d/99-glvd
#!/bin/sh

if [ -x /usr/bin/glvd ]; then
    exec /usr/bin/glvd executive-summary
fi
